# debian/batteriewarner.service
[Unit]
Description=Batteriewarner Daemon

[Service]
Type=simple
ExecStart=/usr/sbin/batteriewarner

[Install]
WantedBy=multi-user.target